# target/linux/mediatek/filogic/base-files/etc/init.d/pwm-fan
#!/bin/sh /etc/rc.common

START=96
 
start() {
	echo "pwmfan started"
	/usr/bin/start-pwmfan.sh > /dev/null&
}